# server/systemd/p2ptalk-updater.timer
[Unit]
Description=p2p-talk auto-update timer (polls for updates)

[Timer]
# First check shortly after boot, then every 15 minutes.
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
